# backups/minecraftbackup.timer
[Unit]
Description=Timer for minecraft server backup

[Timer]
OnUnitActiveSec=172800s
OnBootSec=10s

[Install]
WantedBy=timers.target